# Prelekcja_nr_1/r_httpd_app/templates/iptables.service
[Unit]
After=network.target

[Service]
Type=fork
ExecStart=/sbin/iptables-restore < /etc/iptables/rules.v4

[Install]
WantedBy=multi-user.target

